# tanlov.socket
[Unit]
Description=Tanlov socket

[Socket]
ListenStream=/run/tanlov.sock

[Install]
WantedBy=sockets.target
